# proxmox/systemd/hdd-fanwall-control.service
[Unit]
Description=HDD Fanwall Control

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hdd_fanwall_control.sh